# playbooks/fs-root/default/lib/systemd/system/water_heater_on.timer
[Unit]
Description=water_heater_on.service job

[Timer]
OnCalendar=*-*-* 02:30:00
Unit=water_heater_on.service

[Install]
WantedBy=timers.target
