# fs-root/default/lib/systemd/system/watchdog.service
[Unit]
Description=Trigger watchdog to prove network is still online

[Service]
Type=oneshot
ExecStart=/usr/bin/curl http://hall.local/watchdog/1
