# fs-root/default/lib/systemd/system/water_heater_off.timer
[Unit]
Description=water_heater_off.service job

[Timer]
OnCalendar=*-*-* 07:00:00
Unit=water_heater_off.service

[Install]
WantedBy=timers.target
