# fs-root/default/lib/systemd/system/water_heater_off.service
[Unit]
Description=Disable water heater

[Service]
Type=oneshot
ExecStart=/usr/bin/curl http://hall.local/disable_water_heater/1
